# init/clockwork.conf.erb
description "Clockwork configuration."

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 20

script
  PIDFILE=<%= app_path %>/shared/pids/clockwork.pid
  echo $$ > $PIDFILE
  chown <%= user %>:<%= user %> $PIDFILE
  exec su -c "cd <%= app_path %>/current && RAILS_ENV=production /usr/local/bin/ruby -S bundle exec clockwork <%= app_path %>/current/config/clock.rb -P $PIDFILE >> <%= app_path %>/shared/log/clockwork.log 2>&1" <%= user %>
end script