# logio/logio_server.init.conf.erb
# logio - ensure log.io server starts on boot

description     "log.io server startup script"

start on runlevel [2345]
stop on runlevel [!2345]

exec log.io server start